# ansible/systemd/rentbot-scraper.timer
[Unit]
Description=Run database backup daily

[Timer]
OnCalendar=00,04,08,12,16,20:08:08
Persistent=false

[Install]
WantedBy=timers.target
