# ansible/systemd/rentbot-scraper.service
[Unit]
Description=rentbot run spiders

[Service]
Type=oneshot
ExecStart=/srv/rentbot/run_spiders
WorkingDirectory=/srv/rentbot/
